# ansible/roles/spark/templates/upstart/spark.conf
description "spark-master"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on shutdown

# Respawn the process on unexpected termination
respawn

# respawn the job up to 7 times within a 5 second period.
# If the job exceeds these values, it will be stopped and marked as failed.
respawn limit 7 5

pre-start exec {{ spark_lib_dir }}/sbin/start-master.sh
post-stop exec {{ spark_lib_dir }}/sbin/stop-master.sh
